# ssh/root/etc/s6-overlay/s6-rc.d/init-cloud-run-ssh-server-config/up
/etc/s6-overlay/s6-rc.d/init-cloud-run-ssh-server-config/run
